# scripts/yt-dlp-update.timer
[Unit]
Description=Weekly yt-dlp update check

[Timer]
# Run every Sunday at 4 AM (quiet hours)
OnCalendar=Sun *-*-* 04:00:00
# Run immediately if we missed a scheduled run (e.g., system was off)
Persistent=true
# Add random delay up to 30 min to avoid thundering herd
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
